# config_mgt/files/crontab
*/5 * * * * ruby /home/user/modules/config_mgt/files/installer.rb
